# supervisord.conf
[supervisord]
nodaemon=true

[program:proxy]
command=mitmdump --conf /srv/banditio.proxy/mitmdump.conf
loglevel=debug
