# supervisord.conf
[supervisord]
nodaemon=true

[program:proxy]
command=mitmdump --conf /srv/banditio.proxy/mitmdump.conf

[program:engine]
command=python /srv/banditio.engine/engine.py
